# packages/xorg-startx-systemd/startx.service
[Unit]
Description=StartX service

[Service]
ExecStart=/usr/bin/startx
Restart=on-abort

[Install]
Alias=display-manager.service
